# quadlet/candela-web.container
[Unit]
Description=Candela web UI (FastAPI + Uvicorn)
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/candela:latest
ContainerName=candela-web
EnvironmentFile=/etc/candela/candela.env
PublishPort=8000:8000

Exec=uv run uvicorn candela.main:app --host 0.0.0.0 --port 8000

[Service]
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target default.target
